# deploy/systemd/hermes-backup.timer
[Unit]
Description=Hermes Backup Timer (daily at 02:30)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
