# deploy/systemd/hermes-update.service
[Unit]
Description=Hermes Data Update Service
After=network-online.target

[Service]
Type=oneshot
User=hermes
Group=hermes
WorkingDirectory=/opt/hermes/sembako-dashboard
EnvironmentFile=/etc/hermes-dashboard/hermes-dashboard.env
ExecStart=/opt/hermes/sembako-dashboard/scripts/run_all_updates.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
