# deploy/systemd/hermes-update.timer
[Unit]
Description=Hermes Data Update Timer (every 2 hours)

[Timer]
OnCalendar=*:0/2
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
